# angry-ip-scanner/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start NGINX service
# ==============================================================================

bashio::log.info "Starting Nginx..."
exec nginx
