# cups/rootfs/etc/services.d/dbus/run
#!/bin/sh
exec dbus-daemon --system --nofork