# toogoodtogo-ha-mqtt-bridge/rootfs/etc/services.d/toogoodtogo-ha-mqtt-bridge/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Strating Bridge..."
exec python /app/toogoodtogo_ha_mqtt_bridge/main.py