# cups-dev/rootfs/etc/s6-overlay/s6-rc.d/initialization/up
#!/usr/bin/env bashio

mkdir -p /var/run/dbus

